# config_files/rc.local
#!/bin/sh -e

(sleep 2; /opt/lushroom/start.sh)&

exit 0
